# ibus-voiceinputd.service.in
[Unit]
Description=IBus Voice Input Daemon (Whisper)
Documentation=https://github.com/yuge42/ibus-sample
After=sound.target

[Service]
Type=simple
ExecStart=${HOME}/.local/share/ibus-voiceinputd-venv/bin/python ${HOME}/.local/bin/ibus-voiceinputd
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
# Whisper model loading and inference can be memory intensive.
# Adjust these limits based on your hardware and model size:
# - base/small models: 2G memory may suffice
# - medium model: 4G recommended
# - large models: 8G or more may be needed
# CPUQuota allows using up to 2 CPU cores for parallel processing
MemoryMax=4G
CPUQuota=200%

[Install]
WantedBy=default.target
